# root/etc/s6-overlay/s6-rc.d/log-prepare/up
if { mkdir -p /var/log/myapp }
if { chown root:root /var/log/myapp }
chmod 02755 /var/log/myapp