# root/etc/s6-overlay/s6-rc.d/myapp/run
#!/bin/sh
exec 2>&1
exec mydaemon-in-the-foreground-and-logging-to-stderr